# supervisor/apache2.conf
[program:apache2]
command=/opt/start-apache.sh
autostart=true
autorestart=true

